# deploy/vk-spambot.service
[Unit]
Description=VK SpamBot Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=spambot
WorkingDirectory=/opt/vk-spambot
ExecStart=/usr/bin/python3 -m bot
Restart=always
RestartSec=5
EnvironmentFile=/etc/default/vk-spambot

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict

[Install]
WantedBy=multi-user.target 